# example/openrc#privacy-exposer
#!/sbin/openrc-run

command=/usr/local/bin/privacy-exposer
pidfile="/run/privacy-exposer.pid"
command_args="-l5 -r /var/lib/privacy-exposer/rules.conf -p $pidfile"

extra_started_commands=reload

depend() {
	need net.lo
}

reload() {
	ebegin "Reloading rule set"
	kill -s HUP $(<$pidfile)
	eend $?
}
